# var/lib/dpkg/info/prelink.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/prelink
/etc/default
/etc/default/prelink
/etc/prelink.conf
/usr
/usr/sbin
/usr/sbin/prelink
/usr/sbin/prelink.bin
/usr/share
/usr/share/doc
/usr/share/doc/prelink
/usr/share/doc/prelink/changelog.gz
/usr/share/doc/prelink/README.Debian
/usr/share/doc/prelink/TODO
/usr/share/doc/prelink/copyright
/usr/share/doc/prelink/prelink.pdf.gz
/usr/share/doc/prelink/NEWS.Debian.gz
/usr/share/doc/prelink/changelog.Debian.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/prelink.8.gz
/usr/share/man/man8/prelink.bin.8.gz
