# var/lib/dpkg/info/prelink.conffiles
/etc/cron.daily/prelink
/etc/default/prelink
/etc/prelink.conf
